# general_infos/init_bootimg/inittab
::sysinit:/etc/rc.local
::respawn:/bin/cttyhack /bin/sh
